# config/launchd/com.patabrava.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.patabrava.telegram-bot</string>

  <key>WorkingDirectory</key>
  <string>/Users/highlevel/ScrapperPatabrava</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
    <string>main.py</string>
    <string>telegram-bot</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/highlevel/ScrapperPatabrava/logs/bot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/highlevel/ScrapperPatabrava/logs/bot.log</string>
</dict>
</plist>
